# systemd/azure-ddns.timer
[Unit]
Description=azure-ddns: periodic re-check

[Timer]
# 2 min after boot, then every 10 min. Cheap when nothing changed (one
# ipify probe per enabled stack, no Azure traffic), so 10 min is fine.
# Persistent=true catches up if the host was suspended past a tick.
OnBootSec=2min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
