# systemd/azure-ddns.service
[Unit]
Description=azure-ddns: update Azure DNS with this host's current public IP
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/azure-ddns.env

[Service]
Type=oneshot
ExecStart=/usr/bin/azure-ddns
# Env file is read by the script itself (mode 600). Don't EnvironmentFile=
# here — that would log secret values into journald on read failures.

# Hardening — the script only needs network + a few writable paths.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/azure-ddns /run
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
